# root/etc/init.d/xiaoai-mqtt
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99
# STOP=10 # procd handles stop

start_service() {
    # 确保目录存在
    [ ! -d /var/run ] && mkdir -p /var/run
    [ ! -f /var/log/xiaoai-mqtt.log ] && touch /var/log/xiaoai-mqtt.log

    # 清理残留
    rm -f /var/run/xiaoai-mqtt.status

    # 启动进程
    procd_open_instance
    procd_set_param command /usr/bin/lua "/etc/xiaoai-mqtt/mqtt_client.lua"
    procd_set_param respawn
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param env LUA_PATH="/usr/lib/lua/?.lua;/usr/lib/lua/?/init.lua"
    procd_set_param env LUA_CPATH="/usr/lib/lua/?.so"
    procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "xiaoai-mqtt"
}
