# root/etc/init.d/xiaoai-mqtt
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99
STOP=10
DEPENDS="samba4-client etherwake"

start_service() {
    # 记录初始化日志
    local timestamp=$(date '+%Y-%m-%d %H:%M:%S')
    logger -t xiaoai-mqtt "=========================================="
    logger -t xiaoai-mqtt "XiaoAi MQTT 服务初始化开始 - $timestamp"
    logger -t xiaoai-mqtt "=========================================="

    # 清理残留文件（Lua脚本会在启动时处理进程清理）
    rm -f /var/run/mosquitto_sub.pid /tmp/mosquitto_sub.out /var/run/xiaoai-mqtt.pid /var/run/xiaoai-mqtt.status
    logger -t xiaoai-mqtt "清理旧PID和状态文件完成"
    
    # 创建状态文件
    touch /var/run/xiaoai-mqtt.status
    chmod 644 /var/run/xiaoai-mqtt.status
    echo "mqtt_connection=initializing" > /var/run/xiaoai-mqtt.status
    echo "last_action=服务初始化开始" >> /var/run/xiaoai-mqtt.status
    logger -t xiaoai-mqtt "创建状态文件完成"

    # 设置权限（只在需要时）
    [ ! -x "/etc/xiaoai-mqtt/mqtt_client.lua" ] && chmod 755 /etc/xiaoai-mqtt/mqtt_client.lua
    [ ! -x "/etc/xiaoai-mqtt/status.sh" ] && chmod 755 /etc/xiaoai-mqtt/status.sh
    logger -t xiaoai-mqtt "设置脚本权限完成"

    # 创建日志文件
    [ ! -f "/var/log/xiaoai-mqtt.log" ] && touch /var/log/xiaoai-mqtt.log
    chmod 644 /var/log/xiaoai-mqtt.log
    logger -t xiaoai-mqtt "创建日志文件完成"
    
    # 清理旧的日志轮转文件（优化版，使用更高效的方法）
    for i in 5 6 7 8 9 10; do
        [ -f "/var/log/xiaoai-mqtt.log.$i" ] && rm -f "/var/log/xiaoai-mqtt.log.$i"
    done
    logger -t xiaoai-mqtt "清理旧日志文件完成"

    # 注册服务进程
    procd_open_instance
    # 使用绝对路径的lua解释器，并添加调试输出
    procd_set_param command /usr/bin/lua "/etc/xiaoai-mqtt/mqtt_client.lua"
    procd_set_param respawn
    procd_set_param respawn_threshold 3600
    procd_set_param respawn_timeout 5
    procd_set_param respawn_retry 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param pidfile "/var/run/xiaoai-mqtt.pid"
    # 设置环境变量，确保脚本能够正确访问系统资源
    procd_set_param env PATH=/usr/sbin:/usr/bin:/sbin:/bin
    procd_set_param env LUA_PATH="/usr/lib/lua/?.lua;/usr/lib/lua/?/init.lua"
    procd_set_param env LUA_CPATH="/usr/lib/lua/?.so"
    procd_close_instance
    logger -t xiaoai-mqtt "注册服务进程完成"
    
    # 调试：尝试直接运行Lua脚本查看输出
    logger -t xiaoai-mqtt "调试：尝试直接运行Lua脚本..."
    local debug_output=$(/usr/bin/lua /etc/xiaoai-mqtt/mqtt_client.lua 2>&1 | head -c 500)
    if [ -n "$debug_output" ]; then
        logger -t xiaoai-mqtt "调试输出（前500字符）: $debug_output"
    else
        logger -t xiaoai-mqtt "调试：Lua脚本无输出（可能立即退出）"
    fi

    # 添加启动后验证（增加等待时间并改进检查）
    local wait_count=0
    local max_wait=5  # 最大等待5秒
    
    while [ $wait_count -lt $max_wait ]; do
        sleep 1
        if [ -f /var/run/xiaoai-mqtt.pid ]; then
            # 检查PID文件是否有效（包含数字）
            local pid=$(cat /var/run/xiaoai-mqtt.pid 2>/dev/null)
            if [ -n "$pid" ] && [ "$pid" -eq "$pid" ] 2>/dev/null; then
                logger -t xiaoai-mqtt "服务启动成功，PID: $pid"
                echo "last_action=服务启动成功" >> /var/run/xiaoai-mqtt.status
                break
            else
                logger -t xiaoai-mqtt "警告：PID文件存在但内容无效"
            fi
        fi
        
        wait_count=$((wait_count + 1))
        
        if [ $wait_count -eq $max_wait ]; then
            logger -t xiaoai-mqtt "警告：PID文件未生成（等待超时）"
            echo "last_action=服务启动失败" >> /var/run/xiaoai-mqtt.status
            
            # 尝试检查进程是否在运行（即使没有PID文件）
            local running_pid=$(pgrep -f "lua /etc/xiaoai-mqtt/mqtt_client.lua")
            if [ -n "$running_pid" ]; then
                logger -t xiaoai-mqtt "调试：发现运行中的进程 PID: $running_pid"
                # 尝试手动创建PID文件
                echo "$running_pid" > /var/run/xiaoai-mqtt.pid
                chmod 644 /var/run/xiaoai-mqtt.pid
                logger -t xiaoai-mqtt "已手动创建PID文件"
            fi
        fi
    done
    
    # 记录初始化完成日志
    local end_timestamp=$(date '+%Y-%m-%d %H:%M:%S')
    logger -t xiaoai-mqtt "=========================================="
    logger -t xiaoai-mqtt "XiaoAi MQTT 服务初始化完成 - $end_timestamp"
    logger -t xiaoai-mqtt "=========================================="
}

    logger -t xiaoai-mqtt "=========================================="
}

service_triggers() {
	procd_add_reload_trigger "xiaoai-mqtt"
}
